# app/cherrypy.conf
[supervisord]

[program:punctuation-restore-server]
command=gunicorn wsgi --config /app/gunicorn.conf --chdir=/app